# etc/cron.d/mirror.cron
# File: /etc/cron.d/mirror.cron
# Package: mirror
# Purpose: Run mirror right after work hours each day
05	17	*	*	*	root	/usr/share/mirror/mirror-master.sh
